# files/etc/init.d/ipsec
#!/bin/sh /etc/rc.common

START=90
STOP=10


#ipsec.sonicwall=vpn
#ipsec.sonicwall.left='%any'
#ipsec.sonicwall.leftsourceip='192.168.45.7'
#ipsec.sonicwall.leftsubnet='192.168.45.0/24'
#ipsec.sonicwall.leftid='7.7.7.7'
#ipsec.sonicwall.right='203.25.151.203'
#ipsec.sonicwall.rightid='203.25.151.203'
#ipsec.sonicwall.rightsubnet='172.16.10.0/24'
#ipsec.sonicwall.ike='3des-sha1-modp1024'
#ipsec.sonicwall.esp='3des-sha1-modp1024'
#ipsec.sonicwall.psk='123lan100'

_init_ipsec_config(){
	# PS: Set leftsourceip to setup host-to-site vpn
	#_leftsourceip=`uci get ipsec.sonicwall.leftsourceip`
	_leftsubnet=`uci get ipsec.sonicwall.leftsubnet`
	_leftid=`uci get ipsec.sonicwall.leftid`
	_right=`uci get ipsec.sonicwall.right`
	_rightid=`uci get ipsec.sonicwall.rightid`
	_rightsubnet=`uci get ipsec.sonicwall.rightsubnet`
	_psk=`uci get ipsec.sonicwall.psk`
	_ike=`uci get ipsec.sonicwall.ike`
	_esp=`uci get ipsec.sonicwall.esp`

cat << EOT > /etc/ipsec.conf
# Auto generated by ipsec 
config setup
        # strictcrlpolicy=yes
        # uniqueids = no
        nat_traversal=no

conn sonicwall
        lifetime=28800s
        left=$_left
        leftsourceip=$_leftsourceip
        leftsubnet=$_leftsubnet
        leftid=$_leftid
        right=$_right
        rightid=$_rightid
        rightsubnet=$_rightsubnet
        authby=secret
        ike=$_ike
        auth=esp
        esp=$_esp
        pfs=no
        keyexchange=ikev2
        auto=start
EOT

cat << EOT > /etc/ipsec.secrets
# /etc/ipsec.secrets - strongSwan IPsec secrets file
%any $_rightid : PSK "$_psk"
EOT

}

_init_firewall(){
	_leftsubnet=`uci get ipsec.sonicwall.leftsubnet`
	_rightsubnet=`uci get ipsec.sonicwall.rightsubnet`
   	uci del firewall.vpn.subnet
   	uci add_list firewall.vpn.subnet="$_leftsubnet"
   	uci add_list firewall.vpn.subnet="$_rightsubnet"
	uci commit
}

boot() {
	_init_ipsec_config
	enabled=$(uci get ipsec.sonicwall.enabled)
	if [ $enabled -gt 0 ]; then
		ipsec start
	fi
}

start() {
	_init_ipsec_config
	ipsec start
	_init_firewall
}

stop() {
	ipsec stop
}

restart() {
	ipsec restart
}

reload() {
	ipsec update
}
